# cronjob/crontab
* * * * * root /app/pythonScript.sh
# new line character required for cron!